# profile.d/mysql.sh
PATH=$PATH:$HOME/bin:/usr/local/mysql/bin:/usr/local/mysql/lib
